# docker/supervisor/services.ini
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:sanic]
command=/bin/bash -c "pip3.6 install -U /opt/dist/nginx_aws_cognito* && python3.6 -m nginx_aws_cognito.nginx_aws_cognito"
environment=NGINX_AWS_COGNITO_CONFIG="/opt/sanic.yml"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
